# 1-Apps/1J-Opta-LMX/docs/launchd/com.opta.lmx.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.opta.lmx</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Shared/Opta-LMX/venv/bin/python</string>
        <string>-m</string>
        <string>opta_lmx</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Shared/312/Opta/1-Apps/1J-Opta-LMX</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/var/log/opta-lmx/opta-lmx.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/opta-lmx/opta-lmx.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/Shared/Opta-LMX/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/Shared</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Standard</string>
</dict>
</plist>
